# fbsd/etc/rc.local
#!/bin/sh

echo
echo "Starting btsync"
/home/mytoh/local/apps/btsync/btsync
echo "Started btsync"

exit 0
