# fbsd/etc/rc.local

/usr/local/bin/gosh /home/mytoh/local/repo/kirjasto/skripti/kaiutin.scm
